# system_files/common/usr/etc/profile.d/bazzite-neofetch.sh
#!/bin/bash

# Okadora Fastfetch configuration (replacing Bazzite defaults)
alias fastfetch='/usr/bin/fastfetch -c /usr/share/okadora/fastfetch.jsonc'
alias neofetch='/usr/bin/fastfetch -c /usr/share/okadora/fastfetch.jsonc'
alias neowofetch='/usr/bin/fastfetch -c /usr/share/okadora/fastfetch.jsonc'
alias hyfetch='/usr/bin/fastfetch -c /usr/share/okadora/fastfetch.jsonc'