# system_files/usr/lib/systemd/system/force-niri-session.service
[Unit]
Description=Force Niri as default GDM session
Before=gdm.service
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/bin/force-niri-session.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target