# system_files/usr/lib/systemd/system/okadora-firstboot.service
[Unit]
Description=Okadora User Configuration Check
Documentation=https://github.com/RoudineBWT/okadora
After=local-fs.target network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStartPre=/usr/bin/mkdir -p /var/lib/okadora
ExecStart=/usr/libexec/okadora_firstboot_script.sh
StandardOutput=journal
StandardError=journal